# .config/systemd/user/checkmail.service
[Unit]
Description=check mail
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/my-getmail

[Install]
WantedBy=default.target
